# user/22.03.3-rsync-test/files/etc/uci-defaults/91-custom-add-vagrant-user-and-vagrant-insecure-public-key
# -*- mode: sh -*-
# vi: set ft=sh :

set -eu

# https://openwrt.org/docs/guide-user/additional-software/imagebuilder#restricting_root_access
# https://openwrt.org/docs/guide-developer/uci-defaults
# https://openwrt.org/docs/guide-developer/toolchain/use-buildsystem#custom_files

USER_NAME="vagrant"
USER_SSHPUB="ssh-rsa AAAAB3NzaC1yc2EAAAABIwAAAQEA6NF8iallvQVp22WDkTkyrtvp9eWW6A8YVr+kz4TjGYe7gHzIw+niNltGEFHzD8+v1I2YJ6oXevct1YeS0o9HZyN1Q9qgCgzUFtdOKLv6IedplqoPkcmF0aYet2PkEDo3MlTBckFXPITAMzF8dJSIFo9D8HfdOV0IAdx4O7PtixWKn5y2hMNG0zQPyUecp4pzC6kivAIhyfHilFR61RGL+GPXQ2MWZWFYbAGjyiYJnAmCP3NOTd0jMZEnDkbUvxhMmBYSdETk1rRgm+R4LOzFUGaHqHDLKLX+FIPKcF96hrucXzcWyLbIbEgE98OHlnVYCzRdK8jlqm8tehUc9c9WhQ== vagrant insecure public key"
USER_SHELL="/bin/ash"
SUDO_USER="root"
SUDO_GROUP="sudo"
SUDO_GID="27"
groupadd -r -g "${SUDO_GID}" "${SUDO_GROUP}"
useradd -m -G "${SUDO_GROUP}" -s "${USER_SHELL}" "${USER_NAME}"
passwd -l "${SUDO_USER}" # locks the root user
cat << EOI > /etc/sudoers.d/00-custom
%${SUDO_GROUP} ALL=(ALL) ALL
EOI
USER_HOME="$(eval echo ~"${USER_NAME}")"
mkdir -p "${USER_HOME}"/.ssh
chmod 700 "${USER_HOME}"/.ssh
cat << EOI > "${USER_HOME}"/.ssh/authorized_keys
${USER_SSHPUB}
EOI
chmod 600 "${USER_HOME}"/.ssh/authorized_keys
uci set dropbear.@dropbear[0].PasswordAuth="0"
uci set dropbear.@dropbear[0].RootPasswordAuth="0"
uci commit dropbear
/etc/init.d/dropbear restart

exit 0
